# jmwalletd/rootfs/etc/s6-overlay/s6-rc.d/obwatcher/run
#!/command/with-contenv sh
exec /opt/venv/bin/python -m orderbook_watcher.main
