# jmwalletd/rootfs/etc/s6-overlay/s6-rc.d/jmwalletd-svc/run
#!/command/with-contenv sh
exec /opt/venv/bin/jmwalletd --no-tls
